# templates/launchd/chat.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>__LABEL__</string>

  <!-- MUST be /bin/bash -l <script>. Using -lc causes exit 126. -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-l</string>
    <string>__REPO__/tools/tg-bridge/tg-chat.sh</string>
  </array>

  <!-- Checks Telegram every 2 minutes. Exits instantly if nothing is new,
       so it is cheap to leave running all day. Lower to respond faster. -->
  <key>StartInterval</key>
  <integer>120</integer>

  <key>StandardOutPath</key>
  <string>__REPO__/.tg-bridge/logs/chat.out.log</string>
  <key>StandardErrorPath</key>
  <string>__REPO__/.tg-bridge/logs/chat.err.log</string>
</dict>
</plist>
